# system_files/usr/lib/systemd/user/noctalia.service
[Unit]
Description=Noctalia
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=qs -p /usr/share/lynx/noctalia-shell
Restart=on-failure

[Install]
WantedBy=graphical-session.target 
